# supervisord.conf
# supervisord.conf
# supervisord manages nginx and uvicorn inside the container.
# Logs are written to /logs — a persistent Docker volume mounted from the host.
# This means:
#   - Logs survive container restarts and updates
#   - You can read them directly from your Unraid appdata folder
#   - The admin Live Logs panel reads from the same files
#   - fail2ban on the host can watch these files directly
#
# Rotation: each log file is capped at 10 MB with 5 backups (nginx.log.1 … .5)
# so the total log folder never exceeds ~120 MB.

[supervisord]
nodaemon=true
logfile=/logs/supervisord.log
logfile_maxbytes=2MB
logfile_backups=3
pidfile=/tmp/supervisord.pid

# ── nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off; pid /tmp/nginx.pid;"
autostart=true
autorestart=true
stdout_logfile=/logs/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/logs/nginx.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

# ── uvicorn ───────────────────────────────────────────────────────────────────
[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --access-log
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/logs/uvicorn.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/logs/uvicorn.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
