# supervisord.conf
# supervisord.conf
# supervisord is a simple process manager.
# It starts both nginx and uvicorn when the container launches,
# and restarts either one automatically if they crash.
# The container runs as non-root (appuser) — nginx uses port 8080 to avoid
# needing root privileges to bind port 80.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ── nginx (serves the React frontend on port 8080) ────────────────────────────
[program:nginx]
command=nginx -g "daemon off; pid /tmp/nginx.pid;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── uvicorn (runs the FastAPI backend on port 8000) ───────────────────────────
[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
